# root/etc/services.d/syslog-ng/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid tkf /usr/sbin/syslog-ng -F -f /config/conf/syslog-ng.conf
